# rc-bsd.sh.in
#!/bin/sh
# $Id: rc-bsd.sh.in,v 1.2.6.1 2006/12/08 05:05:19 manu Exp $

# PROVIDE: miltergreylist
# REQUIRE: LOGIN
# BEFORE: mail

. /etc/rc.subr

name="miltergreylist"
rcvar=$name
command="@BINDIR@/milter-greylist"
pidfile="/var/run/milter-greylist.pid"
command_args="-P ${pidfile} -u @USER@"
required_files="/etc/mail/greylist.conf"

load_rc_config ${name}
run_rc_command "$1"
